# vnc.conf
[program:vnc]
command=x11vnc -forever -display :1
priority=15
redirect_stderr=true
stdout_logfile_backups=2
stdout_logfile_maxbytes=10MB
